# autodeploy.service
[Unit]
Description=Auto Deploy Bot from GitHub

[Service]
Type=oneshot
ExecStart=/root/autodeploy.sh
